# services.d/nginx/run
#!/bin/sh

echo "Start nginx"
nginx -g "daemon off;"